# cfg/eve-dirt.rc
#!/bin/sh

. /etc/rc.subr

name="eve-dirt"
rcvar=evedirt_enable

command="INSTALLDIR/bin/dirt-daemon.sh"
command_interpreter=sh
${name}_user=RUNUSER

load_rc_config ${name}
run_rc_command "$1"
