# cfg/jobs.cron
# keep log files small...
0 0 * * * INSTALLDIR/bin/trim-logs.sh
